# scripts/com.podcastbrief.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for `podcastbrief serve` — runs the daily scheduler +
  Telegram RAG bot in one long-running process. Auto-starts on user login,
  auto-restarts on crash, logs to ~/Library/Logs/podcastbrief.

  install-launchd.sh substitutes __INSTALL_DIR__ before copying this into
  ~/Library/LaunchAgents/.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.podcastbrief.serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/.venv/bin/gemma-brief</string>
        <string>serve</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__INSTALL_DIR__/.venv/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Wait at least 30s between restarts so we don't thrash on a hard crash loop. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/logs/gemma-brief.out.log</string>

    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/logs/gemma-brief.err.log</string>
</dict>
</plist>
